# provision/ansible/roles/mining/templates/mining.timer.j2
[Unit]
Description=Run Mining daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target